# scripts/aeroairouter-ui.service.template
[Unit]
Description=AeroAIRouter Config UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__INSTALL_DIR__
ExecStart=__NODE__ __INSTALL_DIR__/src/ui/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=AIROUTER_HOME=__AIROUTER_HOME__

[Install]
WantedBy=default.target
